# addons/priority-email-gateway/src/etc/cron.d/priority-email-gateway
# /etc/cron.d/priority-email-gateway
# Update mail system monitoring data

* * * * *       root    /usr/sbin/priority-email-gateway-updateMonitoringData
