# priority-mail-server/src/etc/cron.d/priority-mail-server
# /etc/cron.d/priority-mail-server
# Update mail system monitoring data

* * * * *       root    /usr/sbin/priority-mail-server-updateMonitoringData
