# priority-mail-server/src/DEBIAN/postrm
#!/bin/sh -e
# Post-removal script for priority-mail-server package

update-rc.d -f postfix-priority remove >/dev/null
update-rc.d -f qpsmtpd-priority remove >/dev/null
